# deploy/systemd/km0-idea-processor.timer
[Unit]
Description=Process KM0 user idea queue (fallback timer)

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=30s

[Install]
WantedBy=timers.target
