# deploy/systemd/km0-idea-processor.path
[Unit]
Description=Watch KM0 user idea queue for new submissions

[Path]
PathModified=/var/spool/km0-ideas/incoming
PathExistsGlob=/var/spool/km0-ideas/incoming/*.json
MakeDirectory=no

[Install]
WantedBy=multi-user.target
