# deploy/systemd/km0-idea-processor.timer
[Unit]
Description=Process KM0 user idea queue (fallback timer, daily)

[Timer]
OnBootSec=5min
OnUnitActiveSec=24h
AccuracySec=1h

[Install]
WantedBy=timers.target
